# docker/etc/services.d/nabarr/run
#!/command/with-contenv sh

exec s6-setuidgid abc /app/nabarr/nabarr run
